# etc/cron.daily/contemplate
#!/bin/bash

/bin/bash /home/sourav/Epictetus/passage.sh > /home/sourav/contemplate_out

#cat $(source /home/sourav/Epictetus/passage.sh) > /home/sourav/contemplate_out
